# config_templates/init.d/notebook.conf
# simple uWSGI script

description "notebook instance"
start on runlevel [2345]
stop on runlevel [06]

respawn

exec /www/notebook/venv/bin/uwsgi /www/notebook/uwsgi.ini
